# dist/smech-farm-orchestrator-1.0.0/systemd/smech-farm-api.service
[Unit]
Description=SmechFarm API and Web UI
After=network.target

[Service]
ExecStart=/usr/local/bin/smech-farm-orchestrator serve
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
