# dist/smech-farm-orchestrator-1.0.0/systemd/smech-farm-daemon.service
[Unit]
Description=SmechFarm Orchestrator Daemon
After=network.target

[Service]
ExecStart=/usr/local/bin/smech-farm-orchestrator daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
